# banco-ai.service.example
[Unit]
Description=BANCO AI Backend
After=network.target

[Service]
Type=simple
User=bancoai
WorkingDirectory=/var/www/banco-ai/backend
ExecStart=/usr/bin/java -jar -Dspring.profiles.active=prod /var/www/banco-ai/backend/target/banco-ai-1.0.0.jar
Restart=always
RestartSec=10

# Database Configuration
# IMPORTANTE: Configure as variáveis de ambiente com suas credenciais reais
Environment="DATABASE_URL=jdbc:postgresql://SEU_HOST_POSTGRES:SEU_PORTA/banco_ai"
Environment="DATABASE_USER=seu_usuario"
Environment="DATABASE_PASSWORD=sua_senha"

# JWT Secret (IMPORTANTE: Altere para uma chave forte em produção!)
Environment="JWT_SECRET=sua_chave_secreta_super_forte_aqui"

# Frontend URL (Ajuste conforme seu domínio)
Environment="FRONTEND_URL=https://seu-dominio.com"

# Email Configuration (Opcional - ajuste conforme necessário)
Environment="MAIL_HOST=smtp.gmail.com"
Environment="MAIL_PORT=587"
Environment="MAIL_USERNAME=seu_email@gmail.com"
Environment="MAIL_PASSWORD=sua_senha_app"

# Upload Directory
Environment="UPLOAD_DIR=/var/www/banco-ai/backend/uploads/carros"

[Install]
WantedBy=multi-user.target

